# grafana.conf
[program:grafana]
command=service grafana-server start
